# GPL_AX12v1/Iplatform/packages/opensource/rpcap/filesystem/etc/init.d/rpcap
#!/bin/sh /etc/rc.common

START=99
start()
{
	 rpcapd -n -p 22334 -d
}

stop()
{
	killall rpcapd
}