# GPL_AX12v1/Iplatform/packages/opensource/net-snmp/filesystem/etc/hotplug.d/firewall/50-snmp
#!/bin/sh

. /lib/functions.sh

if [ "$ACTION" = "start" ]; then
    /etc/init.d/snmpd restart
fi
