# GPL_AX12v1/Iplatform/packages/opensource/wifidog/filesystem/etc/hotplug.d/iface/99-wifidog
#!/bin/sh

[ "$DEVICE" == "lo" ] && exit 0


if [ "$INTERFACE" == "lan"  -a "$ACTION" == "ifup" ]; then
	/etc/init.d/wifidog restart &
fi
